# agent/teknabox-agent.service
[Unit]
Description=TekNaBox Agent
Documentation=https://github.com/your-msp/agent
After=network-online.target
Wants=network-online.target
After=NetworkManager-wait-online.service

[Service]
Type=simple
User=teknabox-agent
Group=teknabox-agent

ExecStart=/usr/local/bin/teknabox-agent
WorkingDirectory=/opt/teknabox-agent

# Restart on any failure, with backoff
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Give the agent 10s to shut down cleanly on SIGTERM before SIGKILL
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL

# Environment
Environment=TEKNABOX_CONFIG_PATH=/etc/teknabox-agent/config.json
Environment=TEKNABOX_AGENT_PATH=/usr/local/bin/teknabox-agent

# Logging — stdout goes to journald automatically
StandardOutput=journal
StandardError=journal
SyslogIdentifier=teknabox-agent

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/etc/teknabox-agent /var/log/teknabox-agent /tmp /opt/teknabox-agent
# Protect against OOM killer — agent is more important than most user processes
OOMScoreAdjust=-100
# Allow network access
AmbientCapabilities=CAP_NET_RAW

# Resource limits — tune for your hardware if needed
# 256M is comfortable for a Pi Zero W; raise for x86 servers running heavy tasks
MemoryMax=256M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
